# docker/osx/crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/15 * * * * root /db-backup.sh >> /var/log/cron.log 2>&1
#
