# docker/laravel-queue/laravel-queue-supervisor.conf
[supervisord]
nodaemon=true

[program:laravel-queue]
command = /usr/bin/php /var/www/artisan queue:listen --sleep 5
user = root
autostart = true
